# deploy/roles/cloudsync/cloud-vault-sync.service
[Unit]
Description=Cloud Vault bidirectional sync
RequiresMountsFor=/srv/cloud-vault
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kay
Group=kay
EnvironmentFile=-/etc/default/cloud-vault-sync
ExecStart=/usr/local/bin/cloud-vault-sync.sh
